# infra/remote/frogsmashers-pg-backup.timer
[Unit]
Description=Daily FrogSmashers PostgreSQL backup

[Timer]
OnCalendar=*-*-* 04:30:00 UTC
RandomizedDelaySec=15m
Persistent=true

[Install]
WantedBy=timers.target
